# pkg/build/image/files/nodes/nomad-client/rootfs/etc/systemd/system/cilium-mounts.service
[Unit]
Description=Cilium init for bpf and cgroupv2 fs

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cilium-mounts
