# pkg/build/image/files/nodes/nomad-client/rootfs/etc/systemd/system/consul-dns.service
[Unit]
Description=Configure systemd-resolved to use Consul DNS
Requires=consul.service
After=consul.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/consul-dns

[Install]
WantedBy=multi-user.target
